# config/supervisord-rsyslog.conf
[program:rsyslog]
command = rsyslogd -n
startsecs = 5
stopwaitsecs = 5
